# cron.d/python-job
50 18 * * * /usr/local/bin/python /app/src/main.py /app/.env >> /var/log/cron.log 2>&1
