# docker/bot/cron.d/python-job
45 18 * * * /bin/bash /app/start.sh >> /var/log/cron.log 2>&1
